# fan.service
[Unit]
Description=Control cooling fan based on CPU temperature
After=multi-user.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/fan.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
